# scripts/codeman-web.service
[Unit]
Description=Codeman Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/arkon/default/claudeman
ExecStart=/usr/bin/node dist/index.js web --https
Restart=always
RestartSec=5
KillMode=process
Environment=NODE_ENV=production
Environment=HOME=/home/arkon
Environment=NODE_COMPILE_CACHE=/home/arkon/.codeman/compile-cache

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codeman

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=default.target
